# bin/sshhub/sshhub.service
[Unit]
Description=Setup a secure tunnel to sshhub
After=network.target
User=pi
Group=pi

[Service]
ExecStart=/usr/bin/ssh -i /home/pi/.ssh/rpi4 honzatomek@sshhub.de -R raspberrypi4:22:localhost:22 -N -o ServerAliveInterval=10

# Restart every >2 seconds to avoid StartLimitInterval failure
RestartSec=5
Restart=always

[Install]
WantedBy=multi-user.target
