# bin/sshhub/sshhub.service
[Unit]
Description=Setup a secure tunnel to sshhub
After=network.target
User=[your username]
Group=[your group]

[Service]
ExecStart=/usr/bin/ssh [USERNAME]@sshhub.de -R [COMPUTERNAME]:22:localhost:22 -N -o ServerAliveInterval=10

# Restart every >2 seconds to avoid StartLimitInterval failure
RestartSec=5
Restart=always

[Install]
WantedBy=multi-user.target
